# deploy/arcadelauncher-server.python.service
[Unit]
Description=ArcadeLauncher private library server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=arcade
Group=arcade
EnvironmentFile=/etc/arcadelauncher-server.env
WorkingDirectory=/opt/arcadelauncher-server
ExecStart=/usr/bin/python3 /opt/arcadelauncher-server/arcade_server.py --host ${ARCADE_HOST} --port ${ARCADE_PORT} --library-root ${ARCADE_LIBRARY_ROOT}
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/srv/arcade-library

[Install]
WantedBy=multi-user.target
