# deploy/arcadelauncher-server.service
[Unit]
Description=ArcadeLauncher private library server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=arcade
Group=arcade
EnvironmentFile=/etc/arcadelauncher-server.env
WorkingDirectory=/opt/arcadelauncher-server
ExecStart=/opt/arcadelauncher-server/arcadelauncher-server --host ${ARCADE_HOST} --port ${ARCADE_PORT} --admin-host ${ARCADE_ADMIN_HOST} --admin-port ${ARCADE_ADMIN_PORT} --library-root ${ARCADE_LIBRARY_ROOT}
Restart=on-failure
RestartSec=5
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/srv/arcade-library

[Install]
WantedBy=multi-user.target
